# messanger-websocket.conf
[program:websockets]
command=/usr/bin/php /var/www/html/laravel-messenger-clone/artisan websockets:serve
numprocs=1
autostart=true
autorestart=true
user=root